# src/Domain0.Service/Domain0.Service.service
[Unit]
Description=Domain0.Service

[Service]
Type=notify
ExecStart=/usr/sbin/Domain0.Service

[Install]
WantedBy=multi-user.target